# metar-print.service
[Unit]
Description=Prints METAR Data

[Service]
Type=oneshot
ExecStart=metar-print